# infrastructure/opencode-server/systemd/opencode-server.service
[Unit]
Description=OpenCode AI Coding Agent Server
Documentation=https://opencode.ai/docs
After=docker.service
Requires=docker.service

[Service]
Type=simple
Restart=always
RestartSec=10

# Working directory
WorkingDirectory=/opt/opencode-server

# Environment file
EnvironmentFile=/opt/opencode-server/.env

# Run as non-root user (create this user first)
User=opencode
Group=opencode

# Docker Compose up
ExecStartPre=-/usr/local/bin/docker-compose -f /opt/opencode-server/docker-compose.yml down
ExecStartPre=-/usr/local/bin/docker-compose -f /opt/opencode-server/docker-compose.yml pull

# Start the service
ExecStart=/usr/local/bin/docker-compose -f /opt/opencode-server/docker-compose.yml up --remove-orphans

# Stop the service
ExecStop=/usr/local/bin/docker-compose -f /opt/opencode-server/docker-compose.yml down

# Reload
ExecReload=/usr/local/bin/docker-compose -f /opt/opencode-server/docker-compose.yml restart

# Timeout settings
TimeoutStartSec=5min
TimeoutStopSec=2min

# Resource limits
# CPUQuota=200%
# MemoryLimit=4G

[Install]
WantedBy=multi-user.target

# Security hardening
[Service]
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/opencode-server
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=opencode-server